# collabora/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Collabora
# Runs the Nginx daemon
# ==============================================================================

# Wait for Collabora to become available
bashio::net.wait_for 43003

bashio::log.info "Starting NGinx..."
exec nginx
