# collabora/rootfs/etc/s6-overlay/s6-rc.d/collabora/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: NCollabora8N
# Runs Collabora script
# ==============================================================================

bashio::log.info "Starting Collabora..."

exec /usr/bin/coolwsd --version --use-env-vars --o:sys_template_path=/opt/cool/systemplate --o:child_root_path=/opt/cool/child-roots --o:file_server_root_path=/usr/share/coolwsd --o:logging.color=false --o:stop_on_config_change=true
