# collabora/rootfs/etc/s6-overlay/s6-rc.d/init-collabora/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Collabora
# Collabora init script, runs before any other service
# ==============================================================================

bashio::log.info "Configuring Collabora..."

declare public_url

rm -f /etc/coolwsd/coolwsd.xml

bashio::var.json \
  public_url "$public_url" \
  ssl "^$(bashio::config 'ssl')" \
  log_level "$(bashio::config 'log_level')" \
  | tempio \
    -template /etc/collabora/config.gtpl \
    -out /etc/coolwsd/coolwsd.xml

if [ "$(bashio::config 'log_level')" == "debug" ]; then
  bashio::log.info "$(cat /etc/coolwsd/coolwsd.xml)"
fi
